# rslsync-healthcheck.service
[Unit]
Description=Resilio Sync Health Check

[Service]
Type=oneshot
ExecStart=%RESILIO_BASE_DIR%/check-rslsync-health.sh
